# ci/server/services/larawhale.target
[Unit]
Description=Larawhale's Services
Requires=default.target
After=default.target
AllowIsolate=yes

[Install]
WantedBy=default.target
